# simple-ddns.service
[Unit]
Description=IPv6 DDNS Monitor and Update Service
After=network-online.target
Wants=network-online.target

[Service]
# Replace 'username' with your actual Linux username
User=username
Group=username

# Path to your script and the working directory
WorkingDirectory=/home/username/scripts
ExecStart=/usr/bin/python3 /home/username/scripts/ddns_monitor.py

Restart=always
RestartSec=120

[Install]
WantedBy=multi-user.target

